# exam/roles/exporters/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=node_exporter
Restart=always
RestartSec=3
ExecStart=/opt/node_exporter/node_exporter

[Install]
WantedBy=multi-user.target
