# lab8/roles/bind/templates/bind_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User=bind
ExecStart=/usr/local/bin/bind_exporter

[Install]
WantedBy=multi-user.target
