# lab9/roles/grafana/templates/grafana-server.service.j2
[Unit]
Description=Grafana-Server

#
[Service]
User=grafana
ExecStart=/opt/grafana/bin/grafana-server --homepath="/opt/grafana/"

[Install]
WantedBy=multi-user.target
