# lab7/roles/exporters/templates/apache_exporter.service.j2
[Unit]
Description=Apache Exporter

[Service]
User=apache_exporter
ExecStart=/opt/apache_exporter/apache_exporter

[Install]
WantedBy=multi-user.target
